# net-watcher.service
[Unit]
Description=Net Watcher - Secure Network Traffic Recorder
Documentation=https://github.com/abja/net-watcher
After=network.target auditd.service
Wants=network.target

[Service]
Type=simple
# Security: Run as dedicated unprivileged user
User=netmon
Group=netmon

# Security: Grant only the packet capture capability
CapabilityBoundingSet=CAP_NET_RAW
AmbientCapabilities=CAP_NET_RAW

# Security: Systemd sandboxing - restrict file system access
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/net-watcher

# Security: Additional sandboxing protections
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
ProtectKernelModules=true
RestrictRealtime=true
RestrictSUIDSGID=true
RemoveIPC=true
NoNewPrivileges=true

# Security: Network restrictions
# Allow only network access for packet capture
# Note: This doesn't restrict the CAP_NET_RAW capability
RestrictAddressFamilies=AF_INET AF_INET6 AF_PACKET
SocketBindDeny=any

# Security: Device access restrictions
PrivateDevices=true
DeviceAllow=/dev/null rw
DevicePolicy=strict

# Security: System call filtering
SystemCallArchitectures=native
SystemCallFilter=@system-service @basic-io @network-io

# Process management
ExecStart=/usr/local/bin/net-watcher serve
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StartLimitInterval=0
StartLimitBurst=0

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=net-watcher

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Working directory
WorkingDirectory=/var/lib/net-watcher

# Environment
Environment=NETWATCHER_DB=/var/lib/net-watcher/dns.sqlite
Environment=NETWATCHER_RETENTION=90
Environment=NETWATCHER_BATCH_SIZE=100

[Install]
WantedBy=multi-user.target